# usr/local/share/qubsd/templates/net-tor/rw/etc/rc.local
#!/bin/sh

# Weird permissions issue has to be resovled for tor to start
chgrp _tor /dev/pf
chmod g+rw /dev/pf	
